# launchd/com.corn.local-ai-gemma-gguf-gateway.plist
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.corn.local-ai-gemma-gguf-gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/corn/.local/bin/local-ai-on-demand-lane-gateway.py</string>
        <string>--lane</string>
        <string>gemma-gguf-coding-fallback</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LOCAL_AI_PLATFORM_ROOT</key>
        <string>/Users/corn/Documents/Corn_Automation/macos-wireless-autoswitch</string>
        <key>LOCAL_AI_RESIDENCY_POLICY</key>
        <string>/Users/corn/.config/local-ai/residency-policy.json</string>
        <key>LOCAL_AI_GATEWAY_LOG</key>
        <string>/Users/corn/.local/state/local-ai-gemma-gguf-gateway.log</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/corn/.local/state/local-ai-gemma-gguf-gateway.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/corn/.local/state/local-ai-gemma-gguf-gateway.launchd.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
