# launchd/com.boneman.linkedin-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.boneman.linkedin-mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/corn/.local/bin/linkedin-scraper-mcp</string>
    <string>--transport</string>
    <string>streamable-http</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>3000</string>
    <string>--path</string>
    <string>/mcp</string>
    <string>--log-level</string>
    <string>WARNING</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/corn/Library/Application Support/BonemanAgentPlatform/logs/linkedin-mcp.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/corn/Library/Application Support/BonemanAgentPlatform/logs/linkedin-mcp.err.log</string>
</dict>
</plist>
